# packaging/webdesktopd.openrc
#!/sbin/openrc-run

description="webdesktopd – persistent web desktop daemon"
command="/usr/bin/webdesktopd"
command_user="webdesktopd:webdesktopd"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_logger="logger -t webdesktopd -p daemon.info"
error_logger="logger -t webdesktopd -p daemon.err"

[ -f /etc/webdesktopd/env ] && . /etc/webdesktopd/env

depend() {
    need net
    after firewall
}
